# supervisord/brewery-tasks.conf
[program:brewery-tasks]
command=python3 manage.py process_tasks --duration 3600
directory=/code
autorestart=true
stopasgroup=true
