# supervisord/brewery.conf
[program:brewery]
command=python3 manage.py runserver 0.0.0.0:8000
directory=/code
autorestart=true
stopasgroup=true
